# ansible/templates/supervisor/supervisor-worker.conf
[program:supervisor-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --queue=default
autostart=true
autorestart=true
numprocs=2
user=root
redirect_stderr=true
